# soal3/crontab
5 6-23/8 * * 0-5 ls -lt /home/devi/bash soal3.sh
